# deploy/prod/boatlab-certbot-renew.service
[Unit]
Description=Boatlab 운영 TLS 인증서 갱신
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/boatlab/prod/release.env
ExecStart=/opt/boatlab/prod/renew-certificate.sh
TimeoutStartSec=30min
